# php-fpm/php-image/crontab/crontab
* * * * * /bin/echo "EveryMin" >> /var/log/cron.log 2>&1
